# cloudflare_tunnel/rootfs/etc/services.d/cloudflare_tunnel/run
#!/usr/bin/with-contenv bashio

if bashio::config.has_value 'url'; then
  URL=$(bashio::config "url")
  exec cloudflared tunnel --url "$URL"
else
  exec cloudflared tunnel --config=/data/options.json run
fi
