# alexa_proxy/rootfs/etc/services.d/alexa_proxy/run
#!/usr/bin/with-contenv bashio

exec node --trace-warnings /index.js
