# cloudflare_tunnel/rootfs/etc/s6-overlay/s6-rc.d/init-log/up
/bin/bash -c /etc/s6-overlay/scripts/init-log.sh