# cloudflare_tunnel/rootfs/etc/services.d/cloudflare-tunnel/run
#!/usr/bin/with-contenv bashio

if bashio::config.has_value 'url'; then
  cloudflared tunnel --config /data/options.yaml
else
  cloudflared tunnel --config /data/options.yaml run
fi

sleep 10
