# cloudflare_tunnel/rootfs/etc/services.d/cloudflare-tunnel/run
#!/usr/bin/with-contenv bashio

cloudflared tunnel --config /data/options.yml run

sleep 10
