# whatsapp/rootfs/etc/services.d/whatsapp/run
#!/usr/bin/with-contenv bashio

exec su-exec pptruser node --trace-warnings /home/pptruser/app/index.js
