# cloudflare_tunnel/rootfs/etc/s6-overlay/s6-rc.d/init-config/up
/bin/bash -c /etc/s6-overlay/scripts/init-config.sh